# crontab/system_temperature_worker.sh
#!/bin/bash

cd ../workers
python system_temp.py
