# crontab/system_temperature_worker.sh
#!/bin/bash
python ../workers/system_temp.py
